# .config/systemd/user/dockd.service
[Unit]
Description=Restore screen layout after dock/undock
[Service]
ExecStart=/usr/bin/dockd --daemon
Restart=always
[Install]
WantedBy=default.target
