# .config/systemd/user/learnpass.timer
[Unit]
Description=Learn passwords

[Timer]
OnCalendar=*-*-*/1 00:00:00

[Install]
WantedBy=default.target
