# .config/systemd/user/test.service
[Service]
ExecStart=/usr/bin/true
Type=oneshot
ExecStartPost=/usr/bin/echo test
