# .config/systemd/user/tmpfiles.service
[Service]
Type=oneshot
ExecStart=systemd-tmpfiles --user --clean --create --remove
